# etc/systemd/system/7c-config-watcher.path
# See the `7c-config-watcher.service` file for notes on the role of this file.
[Path]
PathModified=/opt/7c/panel.conf

[Install]
WantedBy=multi-user.target
