# etc/systemd/system/7c-config-watcher.path
[Path]
PathModified=/opt/7c/panel.conf

[Install]
WantedBy=multi-user.target
